# supervisord.conf
[supervisord]
pidfile = /tmp/supervisord-api.pid

# Change if you need to debug
#loglevel = DEBG

logfile = ~/logs/blockchaininfo_supervisord-crawler.log

[unix_http_server]
file = %(here)s/supervisord.sock

[supervisorctl]
serverurl = unix://%(here)s/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[inet_http_server]
port = *:9010
username = saulomrr
password = larica13

[program:tag]
command = scrapy crawl %(program_name)s -s JOBDIR=crawls/%(program_name)s-%(process_num)s -o %(program_name)s.json -t json
process_name = %(program_name)s/%(process_num)s
numprocs = 1
autostart = false
autorestart = false